# scripts/kidbox.service
[Unit]
Description=KidBox Launcher
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/famulare/git/KidBox
Environment=UV_CACHE_DIR=/tmp/uv-cache
Environment=PYGAME_HIDE_SUPPORT_PROMPT=1
ExecStart=/home/famulare/git/KidBox/.venv/bin/python -m kidbox.launcher
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
